# config/rc.local
pfctl -a redirect-anchor -f /etc/pf.jail 
